# bin/metanome.conf
description "metanome"
author  "Dyn Inc"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

exec /usr/local/bin/node /usr/local/dyn/metanome/bin/daemon.js >> /dev/null 2>&1
